# badminton-booking/badminton-booking.service
[Unit]
Description=ZJNU Badminton Court Auto-Booking

[Service]
Type=oneshot
WorkingDirectory=/root/ZJNU-tycg-order
ExecStart=/root/ZJNU-tycg-order/venv/bin/python /root/ZJNU-tycg-order/badminton-booking/main.py --loop
StandardOutput=journal
StandardError=journal