# badminton-booking/badminton-booking.service
[Unit]
Description=ZJNU Badminton Court Auto-Booking

[Service]
Type=oneshot
WorkingDirectory=/root/ZJNU
ExecStart=/root/ZJNU/venv/bin/python /root/ZJNU/badminton-booking/main.py --loop
StandardOutput=journal
StandardError=journal